# systemd/weather-eink.service
[Unit]
Description=Weather E-Ink Display Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/weather_eink_pizero
ExecStart=/home/pi/weather_eink_pizero/venv/bin/python3 /home/pi/weather_eink_pizero/main.py -c /home/pi/weather_eink_pizero/config.yaml

# Restart on failure with delay
Restart=on-failure
RestartSec=30

# Resource limits (important for Pi Zero)
# Limit memory usage to 100MB
MemoryMax=100M
MemoryHigh=80M

# CPU quota: 50% of one core
CPUQuota=50%

# Nice level: lower priority
Nice=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weather-eink

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pi/weather_eink_pizero/cache /home/pi/weather_eink_pizero/logs

[Install]
WantedBy=multi-user.target
